# .docker/etc/service/application/run
#!/bin/sh

exec 2>&1
source /env

cd /application

exec bin/console application:run:$APP_TYPE